# os-services/linux/elasticsearch.service
[Unit]
Description=Elasticsearch

[Service]
WorkingDirectory=/home/gary/logaggregation/elasticsearch-8.3.0/bin
ExecStart=/home/gary/logaggregation/elasticsearch-8.3.0/bin/elasticsearch
SyslogIdentifier=elasticsearch
User=gary

[Install]
WantedBy=multi-user.target
